# supervisor.conf
[supervisord]
nodaemon=true
logfile=/opt/etherpad-lite-data/supervisor.log
user=etherpad

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:etherpad]
directory=/opt/etherpad-lite
command=node node_modules/ep_etherpad-lite/node/server.js
user=etherpad
autostart=true
autorestart=true
